# systemd/user/power-monitor.service
[Service]
Environment=STARTUP_WAIT=30s
ExecStart=/home/aura/.local/bin/power_monitor.sh

[Install]
WantedBy=default.target
